# templates/default/IIBService.erb
#! /bin/sh

### BEGIN INIT INFO
# Provides:          IBM Integration Bus Node 
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Manages an IBM Integration Bus Node
# Description:       Manages an IBM Integration Bus Node 
#                    
### END INIT INFO

test -x /opt/mqm/bin/endmqm
test -x /opt/ibm/mqsi/iib-10.0.0.7/server/bin/mqsistart
test -x /opt/ibm/mqsi/iib-10.0.0.7/server/bin/mqsistop

#
# Function that starts the daemon/service
#
do_start()
{
<%= @qmgr_listener_comment %>    sudo su - <%= @mq_user_name %> -c 'runmqlsr -m <%= @qmgr_name %> -t TCP -p <%= @qmgr_listener_port %> &'
	sudo su - <%= @user_name %> -c 'mqsistart <%= @node_name %> '
#	return "0"
	return $?
}

#
# Function that stops the daemon/service
#
do_stop()
{
    sudo su - <%= @user_name %> -c 'mqsistop <%= @node_name %> '
<%= @qmgr_listener_comment %>    sudo su - <%= @mq_user_name %> -c 'endmqlsr -m <%= @qmgr_name %>  -t TCP -p <%= @qmgr_listener_port %> '
    sudo su - <%= @user_name %> -c 'endmqm -i <%= @qmgr_name %> ' 
#	return "0"
	return $?
}

case "$1" in
  start)
	do_start
	;;
  stop)
	do_stop
	;;
  restart|force-reload)
	$0 stop
    $0 start
	;;      
  reload)
    echo "Warning - reload does nothing"
    exit 3
    ;;
  *)
	echo "Usage: $0 {start|stop|restart|reload}"
    exit 1
    ;;
esac
:
